# cron/crontab
0 3 * * * php jvforum.site/update_forum_list.php
0 5 * * * php jvforum.site/clean_tokens.php
2 5 * * * php jvforum.site/clean_caches.php
